# rootfs/etc/s6-overlay/s6-rc.d/init-mariadb-bootstrap/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

set -o pipefail

stop_mariadb() {
    mariadb-admin --user=root shutdown
}

# Main function
main() {
    exec > >(while read -r line; do echo "[init-mariadb-bootstrap] ${line}"; done) 2>&1

    # Make sure that the provisioning instance is stopped when the script exits
    trap stop_mariadb EXIT SIGINT SIGTERM

    /app/mariadb-bootstrap

    printf "Bootstrap completed\n"
}

main "$@"
